# deploy/yield-daemon.service
[Unit]
Description=Yield Daemon — Autonomous Yield Generation Infrastructure
Documentation=https://github.com/iaintheardofu/workforce
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Binary + config live read-only in /opt/yield-daemon; the writable state dir is
# the systemd-managed StateDirectory (/var/lib/yield-daemon), owned by the
# transient DynamicUser. The daemon's relative state_dir ("runtime/yield_daemon")
# resolves under this WorkingDirectory.
WorkingDirectory=/var/lib/yield-daemon
StateDirectory=yield-daemon
# --dry-run is kept on purpose: remove it ONLY after configuring live
# endpoints/wallets in config.toml and understanding the risk model.
ExecStart=/opt/yield-daemon/yield-daemon --config /opt/yield-daemon/config.toml --dry-run
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

# --- Sandboxing / hardening ---
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
# Resource guard rails
MemoryMax=2G
TasksMax=512

[Install]
WantedBy=multi-user.target
